# systemd/vnc-billi.service
[Unit]
Description=x11vnc — BILLI MT5 GUI 1회 로그인용 (BILLI msg 455, P2-003-R2 Phase 1)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Environment=DISPLAY=:98
Environment=HOME=/root
ExecStart=/usr/bin/x11vnc \
    -display :98 \
    -localhost \
    -rfbauth /root/.x11vnc/passwd \
    -noxdamage \
    -forever \
    -shared \
    -rfbport 5900 \
    -o /var/log/x11vnc.log
# x11vnc 0.9.16 (Ubuntu 24.04) 는 -bg/-nobg 미지원. systemd Type=simple 이
# foreground 동작 자동 보장 (double-fork 안 함) — 별도 옵션 불필요.
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

[Install]
WantedBy=multi-user.target
